# user/config/systemd/user/dzen.service
[Unit]
Description=Dzen Launch

[Service]
Type=exec
ExecStart=%h/config/choose/scripts/start-dzen.sh
