# user/config/systemd/user/shutdown-reminder.service
[Unit]
Description=Shutdown Reminder Timer

[Service]
Type=oneshot
ExecStart=notify-send "TIME TO SLEEP"

[Install]
WantedBy=xorg.target