# user/config/systemd/user/shutdown-reminder.timer
[Unit]
Description=Shutdown Reminder Timer
Requires=shutdown-reminder.service

[Timer]
OnCalendar=*-*-* 22:30:*


[Install]
WantedBy=xorg.target
