# user/config/systemd/user/shutdown-reminder.service
[Unit]
Description=Shutdown Reminder Timer

[Service]
Type=oneshot
ExecStart=%h/scripts/shutdown-reminder.sh

[Install]
WantedBy=xorg.target