# user/config/systemd/user/test.service
[Service]
Type=oneshot
ExecStart=xrdb -query
