# user/config/systemd/user/xorg.target
[Unit]
Description=xorg Target
Documentation=man:xorg(1)
AllowIsolate=yes
