# os/unix/user/.config/systemd/user/albert.service
[Unit]
Description=Albert
BindsTo=xorg.target
After=xorg.target
PartOf=xorg.target
ConditionPathExists=/usr/bin/albert

[Service]
ExecStart=/usr/bin/albert
ExecRestart=/usr/bin/albert restart
ExecStop=/usr/bin/kill $MAINPID

[Install]
WantedBy=xorg.target
