# user/.config/systemd/user/lemonbar.service
[Unit]
Description=Lemonbar Launch

[Service]
Type=exec
ExecStart=%h/config/choose/scripts/start-lemonbar.sh
