# user/.config/systemd/user/xob.socket
[Unit]
Description=Socket xob reads from
ConditionUser=!root

[Socket]
ListenFIFO=%t/xob.sock
SocketMode=0600

[Install]
WantedBy=xorg.target
